# playbooks/roles/supervisord/templates/supervisord.conf.j2
; Sample supervisor config file.
;
; Generated by Ansible deployment script: modifications will be discarded!

{#
Iterate over supervisord_conf dictionary. Top-level keys are sections,
second-level mappings are options.
#}

{# General supersivor config #}
{% for section in supervisord_conf | dict2items %}
[{{  section.key }}]
{%   for option in section.value | dict2items %}
{{     option.key }}={{ option.value }}
{%   endfor %}
{% endfor %}

{# Programs #}
{% for service in services | dict2items %}
[program:{{ service.key }}]
command={{ service.value.prefix }}/{{ service.value.executable }} {{
    service.value.args
        | default([])
        | map('regex_replace', '^(.*)$', '"\\1"')
        | join(" ") }}
{#   Supervisord options from "control" section #}
{%   for option in service.value.control | default({}) | dict2items %}
{{     option.key }}={{ option.value }}
{%   endfor %}

{% endfor %}
{# Groups #}
{% set groups = services
                | dict2items
                | selectattr("value.group", "defined")
                | groupby("value.group") %}
{% for grouper, list in groups %}
[group:{{ grouper }}]
programs={{ list | map(attribute="key") | join(",") }}

{% endfor %}
